# config/searxng/searxng.service.template
[Unit]
Description=SearXNG Metasearch Engine
After=network.target

[Service]
Type=simple
ExecStartPre=-/usr/bin/podman rm -f searxng
ExecStart=/usr/bin/podman run --rm --name searxng \
  -p 127.0.0.1:__PORT__:8080 \
  -v __SETTINGS_PATH__:/etc/searxng/settings.yml:ro,Z \
  -e SEARXNG_SECRET=searxng-local-only-key \
  __IMAGE__
ExecStop=/usr/bin/podman stop searxng
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
